# deploy/agentwiki.service
[Unit]
Description=AgentWiki - Decentralized Knowledge Base
Documentation=https://github.com/daifei0527/agentwiki
After=network.target

[Service]
Type=simple
User=agentwiki
Group=agentwiki
WorkingDirectory=/opt/agentwiki
ExecStart=/usr/local/bin/agentwiki serve --config /etc/agentwiki/config.yaml
Restart=on-failure
RestartSec=5
LimitNOFILE=65536

# 环境变量
Environment=AGENTWIKI_DATA_DIR=/var/lib/agentwiki
Environment=AGENTWIKI_LOG_LEVEL=info

# 安全设置
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/agentwiki /var/log/agentwiki

[Install]
WantedBy=multi-user.target
